# conf/tsar.cron
# cron tsar collect once per minute
MAILTO=""
* * * * * root /usr/local/tsar/tsar --cron > /dev/null 2>&1
